# little-penguin-1/a04/42_keyboard.rules
SUBSYSTEM=="usb", DRIVERS=="usb", ATTRS{bInterfaceProtocol}==1, RUN+="/root/simple_usb_keyboard_module.ko"
